# systemd/agent.service
[Unit]
Description=Claude Agent
After=network.target

[Service]
Type=simple
User=claude-agent
EnvironmentFile=/etc/agent/env
ExecStart=/usr/local/bin/agent
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
