# systemd/agent.service
# Runs as the dedicated `agent` user (home /mnt/nvme/agent, with its own
# .cargo/.rustup/.nvm). /etc/agent.env holds the app config and also sets
# NVM_DIR + BASH_ENV so the non-login shells the agent spawns can `nvm use`.
# A postgres role `agent` (LOGIN, CREATEDB) backs peer-auth + throwaway test DBs.
[Unit]
Description=Git Agent
After=network.target

[Service]
Type=simple
User=agent
EnvironmentFile=/etc/agent.env
ExecStart=/usr/local/bin/agent
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Sandbox: managed sessions run arbitrary commands (unbound models skip
# approval), so confine writes — the whole filesystem is read-only except the
# agent's home (which is also REPO_BASE_PATH, holding every worktree) and a
# private, auto-cleaned tmp.
ProtectSystem=strict
ReadWritePaths=/mnt/nvme/agent
PrivateTmp=true

[Install]
WantedBy=multi-user.target
